# aria2/tracker.service
[Unit]
Description=aria2 bt tracker
[Service]
ExecStart=/home/data/aria2/tracker.sh
