# aria2/tracker.service
[Unit]
Description=aria2 bt tracker
[Service]
ExecStart=/usr/local/aria2/tracker.sh
